# roles/common/templates/ansible-pull.timer.j2
[Unit]
Description=Ansible pull timer to check for repo updates

[Timer]
{{ ansible_pull_timer }}
Persistent=true

[Install]
WantedBy=timers.target
WantedBy=multi-user.target
